# setup/com.claude-session-monitor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-session-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__APP_DIR__/src/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__APP_DIR__</string>

  <!-- PATH に claude / wezterm / gh が含まれるように（agents/飛ぶ/起動/PRボードで使用） -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/claude-session-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-session-monitor.err</string>
</dict>
</plist>
